# recipes-connectivity/most-driver/most-driver/most-driver-setup.service
[Unit]
Description=Configure the MOST INIC after drivers were loaded.

[Service]
ExecStart=/usr/bin/most-driver-setup-channels.sh
