# recipes-connectivity/most-driver/most-driver/most-driver-setup.path
[Unit]
Description=Wait for MOST INIC driver to create mdev entries.

[Path]
PathExists=/sys/devices/virtual/most/mostcore/devices/mdev0
PathExists=/sys/devices/virtual/most/mostcore/devices/mdev1
PathExists=/sys/devices/virtual/most/mostcore/devices/mdev2

[Install]
WantedBy=multi-user.target
